# contrib/runit/redis/finish
#!/bin/sh

code=$1
status=$2
hostname=$(hostname)
if [ $status -ne 0 ]; then
	echo "Redis FAILED on $hostname: $code, $status"
else
 	echo "Redis stopped on $hostname: $code, $status"
fi
