# contrib/runit/hhvm/finish
#!/bin/sh

code=$1
status=$2
hostname=$(hostname)
if [ $status -ne 0 ]; then
	echo "HHVM FAILED on $hostname: $code, $status"
else
 	echo "HHVM stopped on $hostname: $code, $status"
fi
